# datacube.service
[Unit]
Description=Datacube - Data provider service for application launchers
Documentation=https://github.com/hypercube/datacube
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/datacube
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=info

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadOnlyPaths=/usr
ReadWritePaths=%t/datacube

[Install]
WantedBy=default.target
